# add_files/rock-4d_fs/usr/lib/systemd/system/rkaiq_3A.service
[Unit]
Description=Enable Rockchip camera engine rkaiq
DefaultDependencies=no

[Service]
Type=simple
ExecStart=/usr/bin/rkaiq_3A_server
Restart=no
User=root
Group=root

[Install]
WantedBy=sysinit.target
